# debian/bhindex-syncer.upstart

start on (filesystem and net-device-up IFACE!=lo)

console log

pre-start exec chmod -R 777 /var/lib/bhindex
exec sudo -u nobody /usr/bin/bhindex syncer
post-start script
  sleep 1
  chmod -R 777 /var/lib/bhindex
end script

respawn
respawn limit 5 60
